# hardware/system/aaeon/lamassu-machine.conf
exec /opt/launchers/start-lamassu-machine
respawn
start on (local-filesystems and net-device-up IFACE!=lo)
chdir /opt/apps/machine/lamassu-machine
