# hardware/system/up4000/tejo/supervisor/conf.d/lamassu-watchdog.conf
[program:lamassu-watchdog]
command=/usr/bin/node /opt/lamassu-machine/watchdog.js upboard tejo
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/lamassu-watchdog.err.log
stdout_logfile=/var/log/supervisor/lamassu-watchdog.out.log
stdout_logfile_backups=2
stderr_logfile_backups=2
